# scripts/futurebox.service
[Unit]
Description=FutureBox API Server
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=futurebox
WorkingDirectory=/opt/futurebox/server
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=DATA_DIR=/opt/futurebox/data

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/futurebox/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
